# scripts/deb/postinst


systemctl enable monitor-rust
systemctl start monitor-rust

