# scripts/deb/prerm

systemctl stop monitor-rust
systemctl disable monitor-rust

